# docker/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:git_command]
directory=/app
command=/bin/sh -c 'while true; do echo "pulling" && git pull && git submodule update --init --recursive && sleep 30; done'
autostart=true
autorestart=false
startsecs=0
startretries=3
stderr_logfile=/var/log/supervisor/git_stderr.log
stdout_logfile=/var/log/supervisor/git_stdout.log

[program:hugo]
directory=/app
command=hugo serve --bind=0.0.0.0 --port=80 --appendPort=false --baseURL=https://whoismakingnews-3rgcn42vfa-uc.a.run.app
autostart=true
autorestart=true
startsecs=5
startretries=3
stderr_logfile=/var/log/supervisor/hugo_stderr.log
stdout_logfile=/var/log/supervisor/hugo_stdout.log
